# deploy/bizarre-api/celery-worker.service
[Unit]
Description=Celery worker for jojo-ttrpg-platform
After=network.target redis-server.service postgresql.service

[Service]
User=root
Group=root
WorkingDirectory=/opt/bizarre/backend/src
Environment="PATH=/opt/bizarre/.venv/bin"
Environment="DJANGO_SETTINGS_MODULE=app.settings_prod"
# .env loaded by decouple from WorkingDirectory (see gunicorn.service comment).
ExecStart=/opt/bizarre/.venv/bin/celery -A app worker --loglevel=info
Restart=on-failure

[Install]
WantedBy=multi-user.target
